# systemd/mymuxd.service
[Unit]
Description=mymux daemon — tmux control-mode bridge for the mymux workspace
After=default.target

[Service]
Type=simple
# Installed by scripts/install-systemd.sh (copies the release binary here).
ExecStart=%h/.local/bin/mymuxd

# CRITICAL: only signal mymuxd itself on stop/restart. The tmux server it drives
# lives in this cgroup but MUST survive a mymuxd restart — it holds your sessions.
# (Ephemeral pty shells die on their own when mymuxd's pty masters close.)
# The default KillMode=control-group would kill the tmux server too.
KillMode=process

Restart=on-failure
RestartSec=2

# systemd --user gives services a minimal PATH; include the per-user tool dirs
# so language servers (rust-analyzer via rustup, etc.) resolve. The daemon also
# searches these itself, so this is belt-and-braces.
Environment=PATH=%h/.cargo/bin:%h/.local/bin:%h/go/bin:/usr/local/bin:/usr/bin:/bin

# User environment defaults — proxy settings for clusters whose egress needs
# one (https_proxy / MYMUX_PROXY, picked up by package installs & searches),
# MYMUX_* knobs, extra PATH. `-` = optional. mymuxd and mymux-pkg also read
# this file themselves, so dev runs behave the same.
EnvironmentFile=-%h/.config/mymux/env

# Optional overrides:
# Environment=MYMUX_ROOT=%h
# Environment=MYMUX_ADDR=127.0.0.1:8088
# Environment=MYMUX_SOCKET=mymux

[Install]
WantedBy=default.target
